# net-misc/kea/files/kea-dhcp6-initd
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="Kea DHCPv6 server"

command="/usr/sbin/kea-dhcp6"
command_args="-t ${KEA_CONF:=/etc/kea/${SVCNAME}.conf}"
pidfile="/run/${SVCNAME}.pid"

depend() {
	use net
}
