# net-misc/kea/files/kea-initd
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

: ${KEACTRL_CONF:=/etc/kea/keactrl.conf}

description="Kea DHCP services"
extra_started_commands="reload"
description_reload="Reloads the configuration."

depend() {
	need net
	use logger dns
}

start() {
	ebegin "Starting $description"
	/usr/sbin/keactrl start -c $KEACTRL_CONF
	eend $?
}

stop() {
	ebegin "Stopping $description"
	/usr/sbin/keactrl stop -c $KEACTRL_CONF
	eend $?
}

reload() {
	ebegin "Reloading $description configuration"
	/usr/sbin/keactrl reload -c $KEACTRL_CONF
	eend $?
}
