# net-misc/kea/files/kea-dhcp6-initd
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

# set up extra commands
extra_commands="checkconfig"
extra_started_commands="reload"

# set up configuration info
: ${KEA_DHCP6_CONF:=/etc/kea/kea-dhcp6.conf}
: ${KEA_LOGGER_DESTINATION:=/var/log/kea-dhcp6.log}
: ${KEA_LOCKFILE_DIR:=none}
conf_name=$(basename ${KEA_DHCP6_CONF} | cut -f1 -d'.')

# set up commands and daemon info
# we run as user "dhcp" if we can set capabilities
name="Kea DHCPv6 Server"
command="/usr/sbin/kea-dhcp6"
command_args="${KEA_DHCP6_OPTS} -c ${KEA_DHCP6_CONF}"
[[ -f /sbin/setcap ]] && command_user="dhcp:dhcp"
pidfile="${KEA_LOCKFILE_DIR}/${conf_name}.$(basename ${command}).pid"
start_stop_daemon_args="--background
	--env KEA_LOGGER_DESTINATION=\"${KEA_LOGGER_DESTINATION}\"
	--env KEA_LOCKFILE_DIR=\"${KEA_LOCKFILE_DIR}\"
"

depend() {
	need net
	use dns logger
}

start_pre() {
	
	if [ ! -z ${command_user} ]; then
		checkpath --directory --owner "${command_user}" --mode 0755 \
			"${KEA_LOCKFILE_DIR}" /var/lib/kea
		if [ "x${KEA_LOCKFILE_DIR}" != "xnone" ]; then
			checkpath --file --owner "${command_user}" --mode 0660 \
				"${KEA_LOCKFILE_DIR}"/logger_lockfile
		fi
		/sbin/setcap "CAP_NET_RAW=+ep CAP_NET_BIND_SERVICE=+ep" "${command}"
	fi

	checkconfig || return 1
}

checkconfig() {
	ebegin "Checking configuration file ${KEA_DHCP6_CONF}"
	${command} -t "${KEA_DHCP6_CONF}"
	eend $? "Syntax error in ${KEA_DHCP6_CONF}."
}

reload() {
	
	ebegin "Reloading configuration from file ${KEA_DHCP6_CONF}"
	start-stop-daemon --signal HUP \
		--exec "${command}" --pidfile "${pidfile}"
	eend $?
}
